# config/com.briefing.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.briefing.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/kazuyaishida/daily-briefing/scripts/run-briefing.sh</string>
    </array>

    <!-- ログイン時に1回実行 -->
    <key>RunAtLoad</key>
    <true/>

    <!-- ログ出力 -->
    <key>StandardOutPath</key>
    <string>/Users/kazuyaishida/daily-briefing/output/cron.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kazuyaishida/daily-briefing/output/cron-error.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/kazuyaishida/daily-briefing</string>
</dict>
</plist>
