# config/com.briefing.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.briefing.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/kazuyaishida/daily-briefing/src/server.js</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- ログ出力 -->
    <key>StandardOutPath</key>
    <string>/Users/kazuyaishida/daily-briefing/output/server.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kazuyaishida/daily-briefing/output/server-error.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/kazuyaishida/daily-briefing</string>
</dict>
</plist>
